# recipes-iAloy/iAloy/files/iAloy-onboot-loader.service
[Unit]
Description=iAloy onboot loader

[Service]
Type=simple
Restart=no
StandardOutput=syslog
StandardError=syslog
Environment=QT_QPA_PLATFORM=eglfs
ExecStart=/usr/bin/iAloy

[Install]
WantedBy=multi-user.target

